# target/linux/realtek/base-files/etc/init.d/sysled
#!/bin/sh /etc/rc.common

START=99

start() {
        insmod /lib/modules/2.6.30.9/gpio_rtl819x.ko
        echo 0x340c00 > /sys/module/gpio_rtl819x/parameters/rtl819x_mux
        echo 6 > /sys/class/gpio/export
        echo "out" > /sys/class/gpio/gpio6/direction
        echo "0" > /sys/class/gpio/gpio6/value
} 

stop() {
        rmmod gpio_rtl819x
}
